# scripts/itermremote.host-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.itermremote.host-daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Applications/host_daemon.app/Contents/MacOS/host_daemon</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ITERMREMOTE_HEADLESS</key>
    <string>1</string>
  </dict>

  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/itermremote-host-daemon/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/itermremote-host-daemon/stderr.log</string>
</dict>
</plist>

